# UnixLinux/AnsibleCode/roles/cron/templates/etcbackup.timer
[Unit]
Description=Run 5 mins on boot

[Timer]
OnBootSec=5min
Unit=etcbackup.service

[Install]
WantedBy=multi-user.target