# UnixLinux/AnsibleCode/roles/cron/templates/etcbackup.service
[Unit]
Description=Backup of etc

[Service]
Type=simple
ExecStart=/usr/sbin/etc-backup

[Install]
WantedBy=multi-user.target